# desktop/gnome_calendar_widget_auto/awesome_gnome_xsession/.config/autostart/calendar-indicator-autostart.desktop

[Desktop Entry]
Type=Application
Exec=/usr/share/calendar-indicator/calendar-indicator.py
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[es_ES]=calendar-indicator-autostart.desktop
Name=Calendar-Indicator
Comment[es_ES]=Calendar-Indicator official shortcut
Comment=Calendar-Indicator official shortcut
